# xwiki-platform-distribution/xwiki-platform-distribution-debian/xwiki-platform-distribution-debian-tomcat/xwiki-platform-distribution-debian-tomcat9-common/src/deb/resources/etc/systemd/system/tomcat9.service.d/xwiki-tomcat9-systemd.conf
[Service]
# Need to allow Tomcat9 to read and write in XWiki permanent directory
ReadWritePaths=/var/lib/xwiki/data
