# xwiki-platform-distribution/xwiki-platform-distribution-debian/xwiki-platform-distribution-debian-tomcat/xwiki-platform-distribution-debian-tomcat9-pgsql/src/deb/control/postinst
#!/bin/sh

set -e
#set -x

#########################
# Restart Tomcat
#########################

# Restart tomcat9 service (only if it's active)
if ( systemctl -q is-active tomcat9.service || systemctl -q is-enabled tomcat9.service )
then
  deb-systemd-invoke restart tomcat9.service
fi
