# xwiki-platform-distribution/xwiki-platform-distribution-debian/xwiki-platform-distribution-debian-tomcat/xwiki-platform-distribution-debian-tomcat9-common/src/deb/control/postinst
#!/bin/sh

set -e
#set -x

#########################
# Rights
#########################

## Make sure tomcat is able to write in the data folder
chown -R tomcat9.tomcat9 /var/lib/xwiki/data

#########################
# Configuration
#########################

. /usr/share/debconf/confmodule

#########################
# Restart Tomcat
#########################

invoke-rc.d --quiet tomcat9 restart || {
    RESULT=$?
    # Ignore if tomcat8 init script does not exist (yet)
    if [ $RESULT != 100 ]; then
    exit $RESULT
    fi
}
