# xwiki-platform-distribution/xwiki-platform-distribution-debian/xwiki-platform-distribution-debian-tomcat/xwiki-platform-distribution-debian-tomcat11/xwiki-platform-distribution-debian-tomcat11-pgsql/src/deb/control/postinst
#!/bin/sh

set -e
#set -x

#########################
# Restart Tomcat
#########################

# Restart tomcat service (only if it's active)
if ( systemctl -q is-active tomcat[[debian.tomcat.version]].service || systemctl -q is-enabled tomcat[[debian.tomcat.version]].service )
then
  deb-systemd-invoke restart tomcat[[debian.tomcat.version]].service
fi
